# deploy/launchd/com.hawking.doctor.campaign.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hawking.doctor.campaign.supervisor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3.12</string>
        <string>/Users/scammermike/Downloads/hawking/tools/condense/doctor_campaign_supervisor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/scammermike/Downloads/hawking</string>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/scammermike/Downloads/hawking/reports/condense/general_frontier/DOCTOR_CAMPAIGN/supervisor/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/scammermike/Downloads/hawking/reports/condense/general_frontier/DOCTOR_CAMPAIGN/supervisor/launchd.err.log</string>
</dict>
</plist>
